# rclone-systemd/rclone-mount.service
[Unit]
Description=Rclone Mount Service for Google Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=/usr/bin/mkdir -p "/home/hakan/Google Drive"
ExecStart=/usr/bin/rclone mount "Google Drive:" "/home/hakan/Google Drive" \
    --vfs-cache-mode full \
    --vfs-cache-max-size 20G \
    --vfs-cache-max-age 24h \
    --vfs-read-chunk-size 32M \
    --vfs-read-chunk-size-limit 512M \
    --buffer-size 64M \
    --transfers 4 --checkers 4 \
    --config /home/hakan/.config/rclone/rclone.conf \
    --log-file /home/hakan/.cache/rclone/rclone.log \
    --log-level INFO \
    --retries 5 \
    --low-level-retries 10 \
    --poll-interval 1m \
    --daemon-timeout 10m
ExecStop=/bin/fusermount -uz "/home/hakan/Google Drive"
ExecStopPost=/bin/rmdir "/home/hakan/Google Drive"
Restart=always
RestartSec=10
User=hakan
Group=hakan
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
